# supervisor-provisioner.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:mariadb]
command=/usr/sbin/mysqld

[program:php7.0-fpm]
command=/usr/sbin/php-fpm7.0 --nodaemonize
